# deploy/systemd/dormguard-healthcheck.service
[Unit]
Description=DormGuard health check

[Service]
Type=oneshot
ExecStart=/opt/DormGuard/deploy/monitor/dormguard-healthcheck.sh
