# deploy/systemd/dormguard-backend.service
[Unit]
Description=DormGuard Backend API
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/DormGuard/backend
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/DormGuard/backend/venv/bin/python run.py
Restart=always
RestartSec=5
MemoryMax=256M

[Install]
WantedBy=multi-user.target
