# deploy/systemd/dorm-healthcheck.service
[Unit]
Description=Dorm Power Guard health check

[Service]
Type=oneshot
ExecStart=/opt/dorm-power-guard-lite/deploy/monitor/dorm-healthcheck.sh
